# ratbagd/ratbagd.service.in
[Unit]
Description=Daemon to introspect and modify configurable mice

[Service]
Type=dbus
BusName=org.freedesktop.ratbag1
ExecStart=@sbindir@/ratbagd
Restart=on-abort

[Install]
Alias=dbus-org.freedesktop.ratbag1.service
